# ventilairsec2ha/rootfs/etc/services.d/ventilairsec2ha/finish
#!/bin/bash
# Finish script - cleanup when addon stops

echo "Cleaning up Ventilairsec2HA addon..."
exit 0
